# services/vs-code-server.service
[Unit]
Description=VS Code Server

[Service]
ExecStart=RUN_SHELL RUN_SCRIPT
User=USER_NAME

[Install]
WantedBy=multi-user.target
