# services/vs-code-server.service
[Unit]
Description=VS Code Server

[Service]
ExecStart=/bin/bash RUN_SCRIPT
User=USER_NAME

[Install]
WantedBy=multi-user.target
